# volcpi.service
# VolcPi Monitoring Service (systemd unit)
# ES: Servicio systemd para iniciar la estación VolcPi al arranque.
# EN: Systemd unit to start the VolcPi station at boot.
#
# File: /etc/systemd/system/volcpi.service
# Project: VolcPi
# Purpose: Monitoring of rain, seismic, GPS, battery and LEDs on Raspberry Pi.
# ExecStart: /usr/bin/python3 -u /home/pi/Documents/Projects/volc-pi-project/main.py
# WorkingDirectory: /home/pi/Documents/Projects/volc-pi-project
# Devices: GPS + Seismic (serial by-id), GPIO, I2C
# Requirements: python3, pyserial, lgpio, smbus2, pynmea2; groups: dialout,gpio,i2c
#
# Install:
#   sudo cp /home/pi/Documents/Projects/volc-pi-project/volcpi.service /etc/systemd/system/volcpi.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now volcpi.service
# Logs:
#   sudo journalctl -u volcpi.service -f
# Notes: If /dev/serial/by-id is not visible due to PrivateDevices, consider:
#   BindReadOnlyPaths=/dev/serial/by-id  or set PrivateDevices=no
# Last updated: 2025-10-01

[Unit]
# [ES] Asegura que la red esté completamente operativa (DHCP/resolución) antes de iniciar.
# [EN] Ensures the network is fully up (DHCP/resolution) before starting.
Description=Volc-Pi Monitoring Service
After=network-online.target
Wants=network-online.target

[Service]
# [ES] Ejecuta el proceso en primer plano (recomendado para systemd).
# [EN] Run the process in the foreground (recommended for systemd).
Type=simple

# [ES] Usuario y grupo del proceso.
# [EN] Process user and group.
User=pi
Group=pi

# [ES] Directorio de trabajo del servicio (evita errores CHDIR).
# [EN] Service working directory (prevents CHDIR errors).
WorkingDirectory=/home/pi/Documents/Projects/volc-pi-project

# [ES] Salida sin búfer para logs en tiempo real vía journalctl.
# [EN] Unbuffered output for real-time logs via journalctl.
Environment=PYTHONUNBUFFERED=1

# [ES] (Opcional) Ruta de paquetes del sistema; útil en algunas distros para GPIO.
# [EN] (Optional) System site-packages path; useful on some distros for GPIO.
Environment=PYTHONPATH=/usr/lib/python3/dist-packages

# [ES] Comando principal de la aplicación. Ajusta si usas venv.
# [EN] Main application command. Adjust if you use a virtualenv.
ExecStart=/usr/bin/python3 -u /home/pi/Documents/Projects/volc-pi-project/main.py
# [ES] Alternativa con virtualenv (descomenta y ajusta la ruta si corresponde):
# [EN] Alternative using virtualenv (uncomment and adjust if applicable):
# ExecStart=/home/pi/.venvs/rain/bin/python -u /home/pi/Documents/Projects/volc-pi-project/main.py

# [ES] Retraso breve para dar tiempo a dependencias de hardware/OS.
# [EN] Short delay to allow hardware/OS dependencies to settle.
ExecStartPre=/bin/sleep 3

# [ES] Política de reinicio en caso de fallo.
# [EN] Restart policy in case of failure.
Restart=always
RestartSec=5

[Install]
# [ES] Inicia en el objetivo multiusuario (arranque normal del sistema).
# [EN] Start in the multi-user target (normal system boot).
WantedBy=multi-user.target
