# volcpi.service
# Rainfall Monitoring Service (systemd unit)
# [ES] Archivo de unidad para ejecutar la estación de monitoreo de lluvia al iniciar el sistema.
# [EN] Systemd unit to run the rainfall monitoring station at system boot.

[Unit]
# [ES] Asegura que la red esté completamente operativa (DHCP/resolución) antes de iniciar.
# [EN] Ensures the network is fully up (DHCP/resolution) before starting.
After=network-online.target
Wants=network-online.target

[Service]
# [ES] Ejecuta el proceso en primer plano (recomendado para systemd).
# [EN] Run the process in the foreground (recommended for systemd).
Type=simple

# [ES] Usuario y grupo del proceso.
# [EN] Process user and group.
User=pi
Group=pi

# [ES] Directorio de trabajo del servicio (evita errores CHDIR).
# [EN] Service working directory (prevents CHDIR errors).
WorkingDirectory=/home/pi/Documents/Projects/rain-gauge-project

# [ES] Salida sin búfer para logs en tiempo real vía journalctl.
# [EN] Unbuffered output for real-time logs via journalctl.
Environment=PYTHONUNBUFFERED=1

# [ES] (Opcional) Ruta de paquetes del sistema; útil en algunas distros para GPIO.
# [EN] (Optional) System site-packages path; useful on some distros for GPIO.
Environment=PYTHONPATH=/usr/lib/python3/dist-packages

# [ES] Comando principal de la aplicación. Ajusta si usas venv.
# [EN] Main application command. Adjust if you use a virtualenv.
ExecStart=/usr/bin/python3 -u /home/pi/Documents/Projects/rain-gauge-project/main.py
# [ES] Alternativa con virtualenv (descomenta y ajusta la ruta si corresponde):
# [EN] Alternative using virtualenv (uncomment and adjust if applicable):
# ExecStart=/home/pi/.venvs/rain/bin/python -u /home/pi/Documents/Projects/rain-gauge-project/main.py

# [ES] Retraso breve para dar tiempo a dependencias de hardware/OS.
# [EN] Short delay to allow hardware/OS dependencies to settle.
ExecStartPre=/bin/sleep 3

# [ES] Política de reinicio en caso de fallo.
# [EN] Restart policy in case of failure.
Restart=always
RestartSec=5

# [ES] Capabilidades necesarias para acceso a GPIO (según placa/configuración).
# [EN] Capabilities required for GPIO access (depending on board/config).
PermissionsStartOnly=true
CapabilityBoundingSet=CAP_SYS_RAWIO
AmbientCapabilities=CAP_SYS_RAWIO

# [ES] Endurecimiento sugerido (opcional). Descomentar si no rompe tu app.
# [EN] Suggested hardening (optional). Uncomment if it doesn't break your app.
# NoNewPrivileges=true
# ProtectSystem=full
# ProtectHome=true
# PrivateTmp=true
# RestrictSUIDSGID=true

[Install]
# [ES] Inicia en el objetivo multiusuario (arranque normal del sistema).
# [EN] Start in the multi-user target (normal system boot).
WantedBy=multi-user.target
